# etc/systemd/system/gdm.service.d/override.conf
[Service]
ExecStartPre=/bin/cp /var/home/twau/.config/monitors.xml /var/lib/gdm/.config/monitors.xml